# systemd/user/kiwix.service
[Unit]
Description=Kiwix offline content server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/kiwix-serve --port 8080 --library /home/tristan/kiwix/library.xml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
