# systemd/user/monitor-hotplug.service
[Unit]
Description=Monitor hotplug handler
After=graphical-session.target

[Service]
ExecStart=/bin/sh -c 'prev=$(cat /sys/class/drm/*/status 2>/dev/null | sort); while true; do sleep 2; curr=$(cat /sys/class/drm/*/status 2>/dev/null | sort); [ "$curr" != "$prev" ] && %h/Projects/dotfiles/river/monitor-setup.sh; prev="$curr"; done'
Restart=always
RestartSec=1

[Install]
WantedBy=default.target